# allinone-supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
environment=PYTHONPATH="/app"

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ---- API Server ----
[program:api_server]
command=/bin/sh -c "alembic upgrade head && uvicorn onyx.main:app --host 0.0.0.0 --port 8080"
directory=/app
stdout_logfile=/var/log/api_server.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

# ---- Web Server (Next.js) ----
[program:web_server]
command=node /app/web/server.js
directory=/app/web
environment=PORT="3000",HOSTNAME="0.0.0.0",INTERNAL_URL="http://localhost:8080"
stdout_logfile=/var/log/web_server.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

# ---- Celery Workers ----
[program:celery_worker_primary]
command=celery -A onyx.background.celery.versioned_apps.primary worker
    --loglevel=INFO
    --hostname=primary@%%n
    -Q celery
stdout_logfile=/var/log/celery_worker_primary.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_light]
command=celery -A onyx.background.celery.versioned_apps.light worker
    --loglevel=INFO
    --hostname=light@%%n
    -Q vespa_metadata_sync,connector_deletion,doc_permissions_upsert,checkpoint_cleanup,index_attempt_cleanup,opensearch_migration
stdout_logfile=/var/log/celery_worker_light.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_heavy]
command=celery -A onyx.background.celery.versioned_apps.heavy worker
    --loglevel=INFO
    --hostname=heavy@%%n
    -Q connector_pruning,connector_doc_permissions_sync,connector_external_group_sync,csv_generation,sandbox
stdout_logfile=/var/log/celery_worker_heavy.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_docprocessing]
command=celery -A onyx.background.celery.versioned_apps.docprocessing worker
    --loglevel=INFO
    --hostname=docprocessing@%%n
    -Q docprocessing
stdout_logfile=/var/log/celery_worker_docprocessing.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_user_file_processing]
command=celery -A onyx.background.celery.versioned_apps.user_file_processing worker
    --loglevel=INFO
    --hostname=user_file_processing@%%n
    -Q user_file_processing,user_file_project_sync,user_file_delete
stdout_logfile=/var/log/celery_worker_user_file_processing.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_docfetching]
command=celery -A onyx.background.celery.versioned_apps.docfetching worker
    --loglevel=INFO
    --hostname=docfetching@%%n
    -Q connector_doc_fetching
stdout_logfile=/var/log/celery_worker_docfetching.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_worker_monitoring]
command=celery -A onyx.background.celery.versioned_apps.monitoring worker
    --loglevel=INFO
    --hostname=monitoring@%%n
    -Q monitoring
stdout_logfile=/var/log/celery_worker_monitoring.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
autorestart=true
startsecs=10
stopasgroup=true

[program:celery_beat]
command=celery -A onyx.background.celery.versioned_apps.beat beat
stdout_logfile=/var/log/celery_beat.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
startsecs=10
stopasgroup=true

[program:supervisord_watchdog_celery_beat]
command=python onyx/utils/supervisord_watchdog.py
    --conf /etc/supervisor/conf.d/supervisord.conf
    --key "onyx:celery:beat:heartbeat"
    --program celery_beat
stdout_logfile=/var/log/supervisord_watchdog_celery_beat.log
stdout_logfile_maxbytes=16MB
redirect_stderr=true
startsecs=10
stopasgroup=true

# ---- Log Redirect ----
[program:log-redirect-handler]
command=tail -qF
    /var/log/api_server.log
    /var/log/web_server.log
    /var/log/celery_beat.log
    /var/log/celery_worker_primary.log
    /var/log/celery_worker_light.log
    /var/log/celery_worker_heavy.log
    /var/log/celery_worker_docprocessing.log
    /var/log/celery_worker_monitoring.log
    /var/log/celery_worker_user_file_processing.log
    /var/log/celery_worker_docfetching.log
    /var/log/supervisord_watchdog_celery_beat.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart=true
